# compose_llamav2/setup_h100/service/ipoib.service
[Unit]
Description=Setup IP over Infiniband

[Service]
ExecStart=/adminhome/admin/setup_h100/setup_ib.sh

[Install]
WantedBy=multi-user.target
